# workspace.target
[Unit]
Description=Complete workspace
Wants=kdemodules.target

[Install]
Alias=default.target

