# DE/KDE/kdedmodules.target
[Unit]
Description=KDE modules to run outside of KDE session

